# core/obs-lazy-start/config/claude-observability.service
# systemd --user unit for the observability collector (stage 6, optional).
#
# The lazy-start hook already keeps the collector up per session. This unit is
# the level up: it survives logout/reboot, unattended.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp claude-observability.service ~/.config/systemd/user/
#   # set ExecStart to your CURRENT plugin path (see the caveat below), then:
#   systemctl --user daemon-reload
#   systemctl --user enable --now claude-observability
#   loginctl enable-linger "$USER"   # this box is Linger=no; without this it
#                                    # stops at logout
#
# CAVEAT: the server lives under the plugin, whose path changes on every plugin
# update. Either (a) point ExecStart at a stable checkout you keep updated, or
# (b) re-point ExecStart after each plugin update. The lazy-start hook does not
# have this problem (it resolves the path at runtime), so most people just use
# the hook and skip systemd.

[Unit]
Description=Claude Code observability collector
After=default.target

[Service]
Type=simple
# Edit this path to your server.mjs (absolute). Example shown:
ExecStart=/usr/bin/env node --disable-warning=ExperimentalWarning %h/.claude/plugins/claude-hooks/core/observability/server.mjs
ExecStop=/usr/bin/env node %h/.claude/plugins/claude-hooks/core/observability/server.mjs stop
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
